# 7.4/supervisor/etc/services.d/supervisor/run
#!/usr/bin/with-contenv bash

su webuser

exec /usr/bin/supervisord -n -c /etc/supervisor/supervisord.conf
# /usr/bin/supervisord -n -u webuser -c /etc/supervisor/supervisord.conf